# root/etc/services.d/miniod/run
#!/usr/bin/with-contenv sh

exec /app/minio server /data
